# dot_config/systemd/user/1password.service
[Unit]
Description=Password manager and secure wallet
BindsTo=xwayland-satellite.service
After=xwayland-satellite.service

[Service]
Environment=DISPLAY=:0
Environment=ELECTRON_OZONE_PLATFORM_HINT=auto
ExecStart=/opt/1Password/1password --silent %U
Restart=always
