# dot_config/systemd/user/1password.service
[Unit]
Description=Password manager and secure wallet
BindsTo=graphical-session.target
After=graphical-session.target

[Service]
Environment=ELECTRON_OZONE_PLATFORM_HINT=auto
ExecStart=/opt/1Password/1password --silent %U
Restart=always
