# dot_config/systemd/user/1password.service
[Unit]
Description=Password manager and secure wallet
PartOf=graphical-session.target
Requisite=graphical-session.target
After=xwayland-satellite.service

[Service]
Environment=DISPLAY=:0
ExecStart=/opt/1Password/1password --silent %U
Restart=always

[Install]
WantedBy=graphical-session.target
