# dot_config/systemd/user/waybar.service
[Unit]
Description=Highly customizable Wayland bar

[Service]
ExecStart=/usr/bin/waybar
