# dot_config/systemd/user/niri.service.d/override.conf
[Unit]
Wants=desktop-environment.target
